# api/com.rustlings-api.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rustlings-api.sync</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Update this path after running: cargo build --release -->
        <string>/Users/qmac/Documents/GitHub/rustlings_api/target/release/rustlings_api</string>
        <string>sync</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/qmac/Documents/GitHub/rustlings_api</string>

    <!-- Run on the 1st of every month at 03:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/rustlings-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/rustlings-sync.log</string>
</dict>
</plist>
